# ux_prepare
#!/bin/sh /etc/rc.common

START=91

start() {
  # --- 原有网络配置部分：保持不变 ---
  echo "config interface 'loopback'
    option ifname 'lo'
    option proto 'static'
    option ipaddr '127.0.0.1'
    option netmask '255.0.0.0'

config interface 'lan'
    option type 'bridge'
    option ifname 'eth0'
    option proto 'static'
    option ipaddr '$UU_LAN_IPADDR'
    option gateway '$UU_LAN_GATEWAY'
    option netmask '$UU_LAN_NETMASK'
    list dns '$UU_LAN_DNS'" > /etc/config/network
  /etc/init.d/network restart

  uci set dhcp.lan.ignore=1
  uci commit dhcp
  /etc/init.d/dnsmasq restart

  # --- 优化后的安装与启动逻辑 ---

  # 1. 启动 GOST v3 API
  if [ -x "/usr/bin/gost" ]; then
      # 增加 killall 防止重复启动
      killall gost >/dev/null 2>&1
      /usr/bin/gost -api :"$GOSTAPI" >/dev/null 2>&1 &
  fi

  # 2. 网易 UU：增加持久化检测
  if [ ! -f "/usr/sbin/uu/uuplugin_monitor.sh" ]; then
      # 放在括号里后台运行，防止安装过程卡住后面的奇游逻辑
      (
        wget -q https://uurouter.gdl.netease.com/uuplugin-script/openwrt/v1/install.sh -O /tmp/uuinstall.sh
        [ -s /tmp/uuinstall.sh ] && /bin/sh /tmp/uuinstall.sh openwrt $(uname -m)
      ) &
  else
      /usr/sbin/uu/uuplugin_monitor.sh &
  fi
  
  # 3. 奇游持久化检测
  if [ ! -f "/etc/qy/qy_acc.sh" ]; then
      # 异步执行安装，确保整个 start 流程快速结束
      (
        wget -q http://static.qiyou.cn/upgrade/ljb/plugin/qyplug.sh -O /tmp/qyplug.sh && \
        sh /tmp/qyplug.sh start
      ) &
  else
      /etc/qy/qy_acc.sh start &
  fi
}

stop() {
  killall gost uuplugin qy_acc 2>/dev/null
}



